# overlay/etc/init.d/ccapply
#!/sbin/openrc-run

depend() {
    after network-online
    need net
}

name="ccapply"
command="/usr/sbin/ccapply"
