# overlay/etc/init.d/ccapply
#!/sbin/openrc-run

depend() {
    want network-online
}

name="ccapply"
command="/usr/sbin/ccapply"
