# overlay/etc/init.d/cloud-config
#!/sbin/openrc-run

depend() {
    want network-online
    before ccapply
}

name="cloud-config"
command="/usr/sbin/metadata"
