# overlay/etc/init.d/cloud-config
#!/sbin/openrc-run

depend() {
    after network-online
    need net
}

name="cloud-config"
command="/usr/sbin/metadata"
