# overlay/etc/init.d/issue
#!/sbin/openrc-run

depend() {
    after network-online
    need net
}

name="issue"
command="/usr/sbin/update-issue"
