# overlay/etc/init.d/issue
#!/sbin/openrc-run

depend() {
    use network-online
}

name="issue"
command="/usr/sbin/update-issue"
