# overlay/etc/init.d/setupnetwork
#!/sbin/openrc-run

depend() {
    provide net
}

name="setupnetwork"

start(){
	/sbin/setupnetwork
}
