# images/02-rootfs/overlay/etc/init.d/ccapply
#!/sbin/openrc-run

depend() {
    after net-online
    need net
}

name="ccapply"
command="/usr/sbin/ccapply"
