# overlay/etc/profile.d/vim.sh
alias vi=vim
if [ -z "$EDITOR" ]; then
    export EDITOR=vim
fi
